# roles/docker/files/base.image/supervisor.crond.conf
[program:cron]
command = /usr/sbin/crond
autostart=true
autorestart=true
startretries=10
